# config/includes.chroot/etc/xdg/autostart/light-locker.desktop
[Desktop Entry]
Type=Application
Name=Screen Locker
Comment=Lock screen program
Exec=light-locker
NoDisplay=True
Hidden=True
